# files/etc/uci-defaults/54_set-user-name
sed -i 's/duser = "root"/duser = "visor"/' /usr/lib/lua/luci/dispatcher.lua

sed -i 's/root/visor/' /etc/config/rpcd