# files/etc/uci-defaults/55_set-system
if [ "$(uci get system.@system[0].hostname)" != "NetpingGSM" ]; then
uci set system.@system[0].hostname='NetpingGSM'
uci commit system
fi

uci set dhcp.lan.ignore='1'
uci commit dhcp

uci set luci.languages.ru='Русский'
uci set luci.languages.en='English'
uci set luci.main.lang='ru'
uci commit luci