# files/etc/uci-defaults/56_set-network_modem
uci set network.ModemUSB=interface
uci set network.ModemUSB.proto='dhcp'
uci set network.ModemUSB.ifname='usb0'
