# files/etc/uci-defaults/57_set-cli-link
ln -sf /etc/netping/cli/cli.py /usr/bin/netping