# files/etc/hotplug.d/iface/96-dhcp-restore
#!/bin/sh

[ "$ACTION" = ifup -o "$ACTION" = ifupdate ] || exit 0

res=$(ip a |grep 'eth0$')
icount=$(echo "$res"|wc -l)
[ "$icount" -eq 2 ] || exit 0

/etc/init.d/network restart

